# src/bLsched.service.in
[Unit]
Description=bLsched big/LITTLE cpu manager
After=local-fs.target

[Service]
ExecStart=@sbindir@/bLsched -b 0 -b 1 -t 80

[Install]
WantedBy=multi-user.target

